# bin/crontab
*/5 * * * * node /home/content/scripts/createThumbImages.js > /var/log/createThumbImages
